# bms-cockpit/DEBIAN/postinst
#!/bin/bash
set -e

# Create system user/group
if ! id "bms-cockpit" &>/dev/null; then
    adduser --system --group \
            --home /var/lib/bms-cockpit \
            --shell /bin/bash \
            bms-cockpit

    usermod -aG video,input,tty bms-cockpit
fi

# Ensure correct permissions for sudoers drop-in
SUDOERS_FILE="/etc/sudoers.d/bms-cockpit"
chmod 440 "$SUDOERS_FILE"
chown root:root "$SUDOERS_FILE"

# Reload systemd units & enable service
systemctl daemon-reload
systemctl enable xserver.service
systemctl start xserver.service

# Enable and start pigpiod service
systemctl enable pigpiod
systemctl start pigpiod

exit 0
